# target/classes/deb/control/postinst
#! /bin/sh
update-rc.d daemon-test defaults