# target/classes/deb/control/postrm
#! /bin/sh
userdel -r pi-daemon-test
groupdel pi-daemon-test
update-rc.d -f daemon-test remove