# src/main/resources/deb/control/postinst
#! /bin/sh
update-rc.d ${project.artifactId} defaults