# src/main/resources/deb/control/postrm
#! /bin/sh
userdel -r ${deb.user}
groupdel ${deb.user}
update-rc.d -f ${project.artifactId} remove